# calendar_api/rootfs/etc/services.d/calendarapi/run
#!/usr/bin/with-contenv bashio

# ==============================================================================
# Start the alloy service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

CONFIG_PATH=/data/options.json
if bashio::fs.file_exists "/custom/config/calendarapi.yaml"; then
    bashio::log.info "Using custom calendarapi.yaml configuration"
    CONFIG_PATH="/custom/config/calendarapi.yaml"
fi

exec /usr/bin/calendarapi serve --config $CONFIG_PATH